# scrapius.service
[Unit]
Description=Scrapius Facebook Scraper Bot
After=network.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
WorkingDirectory=/root/scrapius
Environment=PATH=/root/scrapius/venv/bin
ExecStart=/root/scrapius/venv/bin/python /root/scrapius/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scrapius

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/root/scrapius

[Install]
WantedBy=multi-user.target 